# container/mcp-server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:openviking-server]
command=openviking-server --config /app/ov.conf.json
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:openviking-mcp]
command=node /app/mcp-server/server.mjs
directory=/app/mcp-server
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
environment=MCP_PORT=4050,OPENVIKING_BASE_URL=http://127.0.0.1:1934,OPENVIKING_ROOT_API_KEY=%(ENV_OPENVIKING_ROOT_API_KEY)s

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
